# root/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/execlineb -P
nginx -c /config/nginx.conf -g "daemon off;"